# src/cron/custom.cron
#!/usr/bin/env bash
if [ -d "/usr/local/cpanel" ] && [ -x "/var/cpanel/addons/maldet/bin/homematch" ]; then
        #cpanel
        HOMEMATCH=`/var/cpanel/addons/maldet/bin/homematch`
        $inspath/maldet -b -r $HOMEMATCH?/?/public_html $scan_days >> /dev/null 2>&1
fi

/var/cpanel/addons/maldet/bin/update
